# installer/launchd/com.foreman.agent.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.foreman.agent</string>

    <key>ProgramArguments</key>
    <array>
      <string>{{ install_dir }}/.venv/bin/foreman</string>
      <string>gateway</string>
    </array>

    <key>WorkingDirectory</key>
    <string>{{ install_dir }}</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>FOREMAN_POLICY_TARGET</key>
      <string>{{ policy_target }}</string>
      {% for k, v in env_inline.items() %}
      <key>{{ k }}</key>
      <string>{{ v }}</string>
      {% endfor %}
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>{{ home }}/Library/Logs/Foreman/foreman.out.log</string>
    <key>StandardErrorPath</key>
    <string>{{ home }}/Library/Logs/Foreman/foreman.err.log</string>

    <key>ProcessType</key>
    <string>Background</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
  </dict>
</plist>
